# files/root/etc/uci-defaults/luci-app-mnnu
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@mnnu[-1]
	add ucitrack mnnu
	set ucitrack.@mnnu[-1].init=mnnu
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
